# assets/systemd/update-recovery-iso.service
[Unit]
Description=Update Ubuntu recovery ISO

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-recovery-iso
